# rust-api.conf
# Upstart Configuration

description     "Rust Web App"
author          "Jake Scott, Richard Downer"

start on (net-device-up
          and local-filesystems
          and runlevel [2345])
          
stop on runlevel [016]

respawn

exec /path/to/rust-api
